# contrib/systemd/datagerry.service
[Unit]
Description=dataGerry - Enterprise grade OpenSource CMDB
After=network.target

[Service]
Type=simple
ExecStart=/opt/datagerry/bin/datagerry -c /opt/datagerry/etc/cmdb.conf -s
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=datagerry.service
